# contrib/shell-sync-server.service
[Unit]
Description=Shell Sync Server
After=network.target

[Service]
Type=simple
ExecStart=INSTALL_DIR/shell-sync serve --foreground
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
